# puppet/modules/redis/templates/redis_supervisor.conf.erb
[program:redis]
command=/usr/local/bin/redis-server /etc/redis.conf
directory=<%= redis_dir %>
user=root
stopsignal=QUIT
