# puppet/modules/mongodb/templates/mongodb_supervisor.conf.erb
[program:mongodb]
command=<%= mongodb_dir %>/bin/mongod --dbpath <%= mongodb_data_dir %> --rest
directory=<%= mongodb_dir %>
user=root
stopsignal=QUIT
