# puppet/modules/syslog/files/mongodb.conf
[program:mongodb]
command=/opt/mongodb/bin/mongod --dbpath /opt/mongodb_data
directory=/opt/mongodb
user=root
stopsignal=QUIT
