# puppet/modules/syslog/templates/opencloud_syslog.conf.erb
[program:opencloud-syslog]
command=/usr/bin/ruby script/rails server
    -e production -p 3001
directory=<%= graylog_dir %>
user=root
stopsignal=QUIT
