# puppet/modules/syslog/templates/opencloud_syslog_server.conf.erb
[program:opencloud-syslog-server]
command=java -jar graylog2-server.jar 
directory=<%= graylog_server_dir %>
user=root
stopsignal=QUIT
